# setup/ycontroller_root/etc/systemd/system/arena-attached@.service
[Unit]
Description=Arena attached

[Service]
User=yuser
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/local/bin/find_arenas_then_run.py
