# setup/ycontroller_root/etc/systemd/system/arena-attached.service
[Service]
Type=oneshot
ExecStart=/usr/local/bin/arena_count.sh
