# setup/ycontroller_root/etc/systemd/system/arena-attached.service
[Service]
Type=simple
TimeoutStartSec=0
ExecStartPre=-/usr/bin/docker stop %n
ExecStart=/usr/local/bin/find_arenas_then_run
