# setup/ycontroller_root/etc/systemd/system/ycontroller.timer
[Unit]
Description="Start ycontroller on boot after delay."

[Timer]
OnBootSec=20s

[Install]
WantedBy=default.target
