# setup/ycontroller_root/etc/systemd/system/ycontroller.service
[Unit]
Description=Find arenas then run

[Service]
Type=oneshot
ExecStart=/usr/local/bin/find_arenas_then_run
