# setup/ycontroller_root/etc/systemd/system/ycontroller.service
[Unit]
Description="Start ycontroller on boot after delay."

[Service]
Type=oneshot
ExecStart=/usr/local/bin/find_arenas_then_run
TimeoutStartSec=0
