# setup/ycontroller_root/etc/systemd/system/ycontroller.timer
[Unit]
Description=ycontroller timer
After=arena-attached.service
Requires=arena-attached.service

[Timer]
OnActiveSec=5s
RemainAfterElapse=False
