# setup/ycontroller_root/etc/systemd/system/ycontroller.service
[Unit]
Description=Finds arenas then runs container

[Service]
Type=oneshot
ExecStart=/usr/local/bin/arena_count.sh
